# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:java-backend]
command=java -jar /app/ruoyi-admin.jar
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/java-backend.log
stderr_logfile=/var/log/supervisor/java-backend.error.log
user=root

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.error.log
user=root